# systemd/user/hypridle.service
[Unit]
Description=Hyprland idle daemon
Documentation=https://github.com/hyprwm/hypridle
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=hypridle
Restart=on-failure

[Install]
WantedBy=graphical-session.target
