# systemd/user/walker.service
[Unit]
Description=Walker app launcher (GApplication service)
Documentation=https://github.com/abenz1267/walker
After=graphical-session.target elephant.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=walker --gapplication-service
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
